# systemd/paradex-camera-agent@.service
[Unit]
Description=Paradex Aravis/GStreamer camera agent for %i
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/opt/paradex
Environment=PYTHONUNBUFFERED=1
# Set this to 1400 only when the camera NIC/switch path cannot use jumbo frames.
Environment=PARADEX_GIGE_PACKET_SIZE=9000
ExecStart=/opt/paradex/.venv/bin/python /opt/paradex/src/camera/server_daemon.py --backend aravis-gstreamer
Restart=on-failure
RestartSec=2
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
